# scripts/amneziawg.service
[Unit]
Description=AmneziaWG kernel module
After=network-pre.target
Before=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
TimeoutStopSec=30

# Check kernel version matches built module
ExecStartPre=/bin/bash -c '\
    EXPECTED=$(cat /data/amneziawg/.kernel-version 2>/dev/null || echo "unknown"); \
    CURRENT=$(uname -r); \
    if [ "$EXPECTED" != "$CURRENT" ]; then \
        echo "Kernel changed: $EXPECTED -> $CURRENT. Module NOT loaded. Rebuild required." | logger -t amneziawg; \
        exit 1; \
    fi'

# Load dependencies
ExecStartPre=/sbin/modprobe udp_tunnel
ExecStartPre=/sbin/modprobe ip6_udp_tunnel
ExecStartPre=/bin/mkdir -p /data/amneziawg/tunnels

# Load AmneziaWG module
ExecStart=/sbin/insmod /data/amneziawg/amneziawg.ko

# Bring up configured interfaces
ExecStartPost=/bin/bash -c '\
    for conf in /data/amneziawg/tunnels/*.conf; do \
        [ -f "$conf" ] || continue; \
        AWG=/data/amneziawg/awg bash /data/amneziawg/awg-quick up "$conf" || true; \
    done'

# Tear down interfaces on stop
ExecStop=/bin/bash -c '\
    for conf in /data/amneziawg/tunnels/*.conf; do \
        [ -f "$conf" ] || continue; \
        AWG=/data/amneziawg/awg bash /data/amneziawg/awg-quick down "$conf" || true; \
    done'
ExecStop=/sbin/rmmod amneziawg

[Install]
WantedBy=multi-user.target
